# home/Library/LaunchAgents/com.user.vscode-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.vscode-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/zsh</string>
        <string>-c</string>
        <string>source ~/.zshrc && vscode-backup && cd ~/.local/share/chezmoi && git add . && git commit -m "Automated VSCode backup - $(date)"</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer> <!-- Sunday -->
        <key>Hour</key>
        <integer>18</integer> <!-- 6 PM -->
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/vscode-backup.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/vscode-backup-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/jarodtaylor</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
